# dot_config/systemd/user/backup.service
[Unit]
Description=Backs up files

[Service]
ExecStart=/usr/bin/rsync %h/notes/todo.txt %h/notes/Passwords.kdbx -CErltm --delete arpi:~/backup
