# dot_config/systemd/user/swww.service
[Unit]
Description=swww daemon

[Service]
Type=simple
ExecStart=/usr/bin/swww-daemon
Restart=always
RestartSec=60

[Install]
WantedBy=sway-session.target
